# systemd/cairn-ai.service
[Unit]
Description=Cairn local AI inference server (llama.cpp)
After=network.target
ConditionPathExists=/opt/cairn/models/gemma-3-1b-it.gguf

[Service]
User=cairn
ExecStart=/usr/local/bin/llama-server \
    --model /opt/cairn/models/gemma-3-1b-it.gguf \
    --port 8081 \
    --host 127.0.0.1 \
    --ctx-size 2048 \
    --threads 4
Restart=always
RestartSec=5
MemoryMax=2G

[Install]
WantedBy=multi-user.target
